# debian/postinst
#!/bin/sh

set -e

update-rc.d heka defaults
