# debian/postrm
#!/bin/sh

USER="heka"
GROUP="$USER"

set -e

if [ "$1" = "purge" ]
then
	update-rc.d heka remove
	deluser $USER
	delgroup $GROUP
	rm -rf /var/cache/hekad
fi
