# debian/postrm
#!/bin/sh

set -e

if [ "$1" = "purge" ]
then
	update-rc.d heka remove
	deluser heka
fi
